# templates/default/systemd/etcd.service.erb
[Unit]
Description=Etcd Application Container Engine
Documentation=https://coreos.com/etcd
After=network.target

[Service]
Type=notify
<% if @config.run_user %>
User=<%= @config.run_user %>
<% end %>
<% if @config.http_proxy %>
Environment="HTTP_PROXY=<%= @config.http_proxy %>"
<% end %>
<% if @config.https_proxy %>
Environment="HTTPS_PROXY=<%= @config.https_proxy %>"
<% end %>
<% if @config.no_proxy %>
Environment="NO_PROXY=<%= @config.no_proxy %>"
<% end %>
ExecStart=<%= @etcd_cmd %>
ExecStartPost=/usr/libexec/<%= @etcd_name %>-wait-ready
Restart=always
MountFlags=slave
LimitNOFILE=1048576
LimitNPROC=1048576
LimitCORE=infinity

[Install]
WantedBy=multi-user.target
