# templates/default/etcd.conf.erb
description "etcd service registry"

start on started networking
stop on shutdown

exec /usr/local/bin/etcd <%= @args %> -data-dir <%= node[:etcd][:state_dir] %> -name <%= node[:hostname] %> <% if node[:etcd][:snapshot] %>-snapshot<% end %>
