# templates/default/etcd.conf.erb
description "etcd service registry"

start on started
stop on shutdown

exec /usr/local/bin/etcd <%= @args %> -d <%= node[:etcd][:state_dir] %> -n <%= node[:fqdn] %>
