# templates/default/etcd.service.erb
[Unit]
Description=etcd

[Service]
User=<%= node[:etcd][:user] %>
PermissionsStartOnly=true
ExecStart=/usr/local/bin/etcd <%= @args %> -data-dir <%= node[:etcd][:state_dir] %>
Restart=always
RestartSec=10s

[Install]
WantedBy=default.target
