# templates/default/etcd.conf.erb
description "etcd service registry"

<%- if node[:etcd][:respawn] %>
respawn
<%- end %>

setuid <%= node[:etcd][:user] %>
setgid <%= node[:etcd][:group] %>

limit nofile <%= node[:etcd][:upstart][:nofile] %> <%= node[:etcd][:upstart][:nofile] %>

start on <%= node[:etcd][:upstart][:start_on] %>
stop on <%= node[:etcd][:upstart][:stop_on] %>

exec /usr/local/bin/etcd <%= @args %> -data-dir <%= node[:etcd][:state_dir] %>
