# templates/default/etcd.conf.erb
description "etcd service registry"

start on <%= node[:etcd][:upstart][:start_on] %>
stop on <%= node[:etcd][:upstart][:stop_on] %>

exec /usr/local/bin/etcd <%= @args %> -data-dir <%= node[:etcd][:state_dir] %>
