# templates/default/etcd.conf.erb
description "etcd service registry"

start on started
stop on shutdown

exec /usr/local/bin/etcd <%= node[:etcd][:args] %> -d <%= node[:etcd][:state_dir] %> <%= node[:etcd][:name_switch] %>  <%= node[:fqdn] %>
